# webtornado/debian/cron.d
* * * * * 	www-data	/usr/share/webtornado/bin/manager.pl